# deploy/armosphera-one.service.tmpl
[Unit]
Description=Armosphera One Claude — sovereign Armenian business OS
After=network.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=__ROOT__
Environment=PORT=__PORT__
Environment=NODE_ENV=production
ExecStart=__NODE__ __ROOT__/server/index.js
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
